# dockerhub/forward/app.ini
[program:app]
command=/usr/local/bin/node /var/app/index.js
autostart=true
startsecs=3
stdout_logfile=/var/log/supervisor/access.log
stdout_logfile_maxbytes=10MB
redirect_stderr=true
stderr_logfile=/var/log/supervisor/error.log
stderr_logfile_maxbytes=10MB
stopsignal=QUIT
